# meter/datel.service
[Unit]
Description=DATEL DLMS Server

[Service]
ExecStart=/var/energy-network-arena/meter/service.sh

[Install]
WantedBy=multi-user.target
